# setup/jarvis-api.service
[Unit]
Description=Jarvis API Server with Starline GPS Monitor
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/geodrive_n8n-agents
Environment="NODE_ENV=production"
EnvironmentFile=/root/geodrive_n8n-agents/.env
ExecStart=/usr/bin/npm run start:api
Restart=always
RestartSec=10
StandardOutput=append:/var/log/jarvis-api.log
StandardError=append:/var/log/jarvis-api-error.log

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Автоматический перезапуск при падении
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
